# config/dashcam-record.service
[Unit]
Description=Dash Car Cam - Recording Service
After=multi-user.target
Wants=dashcam-cleanup.timer

[Service]
Type=simple
User=pi
EnvironmentFile=/etc/dashcam/dashcam.conf
ExecStart=/usr/local/bin/dashcam-record.sh
KillSignal=SIGINT
TimeoutStopSec=15
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
