# config/dashcam-stream.service
[Unit]
Description=Dash Car Cam - RTSP Streaming Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
EnvironmentFile=/etc/dashcam/dashcam.conf
ExecStart=/usr/local/bin/dashcam-stream.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
